# sot/sync/systemd/ilma-sot-sync.timer
[Unit]
Description=ILMA SOT full delta sync every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=5min
Persistent=true
Unit=ilma-sot-sync.service

[Install]
WantedBy=timers.target
